# roles/caddy/templates/tailscale-cert-renew.timer.j2
[Unit]
Description=Daily Tailscale TLS cert renewal

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
