# roles/ollama/templates/ollama.service.j2
[Unit]
Description=Ollama (containerized, GPU)
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker rm -f ollama
ExecStart=/usr/bin/docker run --rm --name ollama \
    --gpus all \
    -p {{ ollama_bind }}:11434 \
    -v {{ ollama_data_volume }}:/root/.ollama \
    {{ ollama_image }}
ExecStop=/usr/bin/docker stop ollama
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
